# server/etc/systemd/unitychan.service
# Systemd service configuration for unitychan
[Unit]
Description=iFun Engine game - unitychan
After=syslog.target network-online.target

[Install]
WantedBy=multi-user.target

[Service]
# If this limit is not enough, you SHOULD look after
# "net.netfilter.nf_conntrack_tcp_timeout_time_wait" (using sysctl)
LimitNOFILE=999999

# It is strongly recommended that one should use uid:gid other than the root
# NOTE: Please update User=, Group=, and ExecStartPre=
User=root
Group=root

# Makes the daemon automatically restart
Type=simple
Restart=always
RestartSec=5s

# create directories, and change permissions (with root privileges)
PermissionsStartOnly=true
ExecStartPre=/usr/bin/mkdir -p /var/log/funapi/unitychan/glog \
  /var/log/funapi/unitychan/activity \
  /var/log/funapi/unitychan/counter \
  /var/crash/funapi/unitychan
# NOTE: change root:root to service's uid:gid
ExecStartPre=/usr/bin/chown -R root:root /var/log/funapi/unitychan/glog \
  /var/log/funapi/unitychan/activity \
  /var/log/funapi/unitychan/counter \
  /var/crash/funapi/unitychan

EnvironmentFile=/etc/default/unitychan
ExecStart=/usr/bin/unitychan-launcher
